# .config/systemd/user/bat-shutdown.service
[Unit]
Description=Shutdown on low battery.

[Service]
Type=oneshot
ExecStart=/home/gautam/bin/bat-shutdown