# .config/systemd/user/checkmail.service
[Unit]
Description=Mailbox synchronization service

[Service]
Type=oneshot
ExecStart=/home/gautam/bin/checkmail

[Install]
WantedBy=default.target