# .config/systemd/user/checkmail.timer
[Unit]
Description=Mailbox synchronization timer

[Timer]
OnBootSec=2m
OnUnitActiveSec=5m
Unit=checkmail.service

[Install]
WantedBy=timers.target