# .config/systemd/user/inkscapefigures.service
[Unit]
Description=Intercept keypresses in Inkscape.

[Service]
Type=simple
RemainAfterExit=no
ExecStart=/usr/bin/python /home/gautam/src/inkscape-shortcut-manager/main.py

[Install]
WantedBy=default.target