# .config/systemd/user/inkscapewatch.service
[Unit]
Description=Watch for new Inkscape files.

[Service]
Type=oneshot
ExecStart=/home/gautam/bin/inkscape-figures watch