# .config/systemd/user/monitor.service
[Unit]
Description=Monitor hotplug

[Service]
Type=simple
RemainAfterExit=no
ExecStart=/home/gautam/bin/dock

[Install]
WantedBy=default.target