# .config/systemd/user/syncmedia.service
[Unit]
Description=Sync media service

[Service]
Type=oneshot
ExecStart=/home/gautam/bin/syncmedia