# .config/systemd/user/syncmedia.timer
[Unit]
Description=Sync media onto external drive

[Timer]
OnBootSec=2m
OnUnitActiveSec=60m
Unit=syncmedia.service

[Install]
WantedBy=timers.target