# .config/systemd/user/tlp-usb.service
[Unit]
Description=Setup tlp usb

[Service]
Type=oneshot
User=root
ExecStar
t=/usr/bin/tlp usb

[Install]
WantedBy=default.target